# whonix_shared/etc/xdg/autostart/whonixcheck.desktop
[Desktop Entry]
Type=Application
Exec=/usr/lib/whonix/output --identifier whonixcheck --icon /usr/share/whonix/icons/whonix.ico --showx --autostart
